# espocrm-docker-9.1.6/pdf-service/pdf-service.service
[Unit]
Description=PDF Generation Service (Localhost Web API)
After=network.target

[Service]
Type=simple
User=pdf-service
Group=pdf-service

# Run PHP built-in server on localhost only
ExecStart=/usr/bin/php -S 127.0.0.1:8888 -t /opt/pdf-service /opt/pdf-service/index.php

# Working directory
WorkingDirectory=/opt/pdf-service

# Restart policy
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/tmp/pdf-service /opt/pdf-service/logs

# Resource limits
MemoryLimit=512M
CPUQuota=50%

# Logging to file (readable by users)
StandardOutput=append:/opt/pdf-service/logs/stdout.log
StandardError=append:/opt/pdf-service/logs/stderr.log

[Install]
WantedBy=multi-user.target